# image/systemctl-enable-linger.service
[Unit]
Description=Enable Linger for Stratux-Radar

[Service]
ExecStart=loginctl enable-linger pi

[Install]
WantedBy=default.target